# reflection.service
[Service]
User=root
ExecStart=/usr/local/bin/reflection

[Install]
WantedBy=multi-user.target
